# server/supervisord/supervisor.d/openrvdas.ini
[program:cached_data_server]
command=/usr/local/bin/python3 server/cached_data_server.py --port 8766 --disk_cache /var/tmp/openrvdas/disk_cache --max_records 8640 -v
directory=/opt/openrvdas
autostart=false
autorestart=true
startretries=3
stderr_logfile=/var/log/openrvdas/cached_data_server.err.log
stdout_logfile=/var/log/openrvdas/cached_data_server.out.log
user=pablo

[program:logger_manager]
command=/usr/local/bin/python3 server/logger_manager.py --database django --no-console --data_server_websocket :8766 -v
directory=/opt/openrvdas
autostart=false
autorestart=true
startretries=3
stderr_logfile=/var/log/openrvdas/logger_manager.err.log
stdout_logfile=/var/log/openrvdas/logger_manager.out.log
user=pablo

[group:openrvdas]
programs=cached_data_server,logger_manager

[program:simulate_nbp_serial]
command=/usr/local/bin/python3 logger/utils/simulate_serial.py --config test/NBP1406/serial_sim_NBP1406.yaml --loop
directory=/opt/openrvdas
autostart=false
autorestart=true
startretries=3
stderr_logfile=/var/log/openrvdas/simulate_serial.err.log
stdout_logfile=/var/log/openrvdas/simulate_serial.out.log
user=pablo

[group:simulate]
programs=simulate_nbp_serial
